# systemd.dofler.service
[Unit]
Description=Dofler Service
After=mariadb.service

[Service]
Restart=always
WorkingDirectory=/opt/WallOfSheep
Environment=NODE_ENV=production
ExecStart=/root/.nvm/current/bin/node /opt/WallOfSheep/server.js

[Install]
WantedBy=multi-user.target
